# scripts/deploy/presenter.service
[Unit]
Description=Presenter - Church Worship Display System
After=network.target

[Service]
Type=simple
User=newlevel
Group=newlevel
WorkingDirectory=/opt/presenter
ExecStart=/opt/presenter/presenter-server
Restart=always
RestartSec=5

# Environment
Environment=PRESENTER_PORT=80
Environment=PRESENTER_DB_URL=sqlite:///opt/presenter/presenter.db
Environment=PRESENTER_LIBRARY_ROOT=/home/newlevel/devel/presenter/presenter-libraries
Environment=RUST_LOG=info,tower_http=debug

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/presenter

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
