# scripts/deploy/presenter-dev.service
[Unit]
Description=Presenter Dev - Church Worship Display System (Development)
# #333 item 2: order start after the DRI render node so the GStreamer
# `va` plugin can probe the device cleanly. After= is ordering-only,
# does not pull the device in. See presenter.service for full notes.
After=network.target dev-dri-renderD128.device

[Service]
Type=simple
User=newlevel
Group=newlevel
WorkingDirectory=/opt/presenter-dev
# #339: gate start on hardware H264 encoder being probeable. Mirrors
# prod presenter.service; see that file for the full incident
# context. Probes EITHER vah264enc (prod N100 Intel) OR nvh264enc
# (dev2 RTX Nvidia). Leading `-` makes the gate non-fatal — if no
# GPU is available at all, presenter still starts for non-NDI
# features and item 6 skips NDI auto-restore.
ExecStartPre=-/usr/bin/timeout 30 sh -c 'until gst-inspect-1.0 --exists vah264enc || gst-inspect-1.0 --exists nvh264enc; do sleep 1; done'
ExecStart=/opt/presenter-dev/presenter-server
Restart=always
RestartSec=5

# Environment
Environment=PRESENTER_PORT=8080
Environment=PRESENTER_DB_URL=sqlite:///opt/presenter-dev/presenter.db
Environment=PRESENTER_LIBRARY_ROOT=/home/newlevel/devel/presenter/presenter-libraries
# #404: stage URL the launcher opens on every Android stage display (see
# presenter.service). Dev points at its own /stage on port 8080.
Environment=PRESENTER_ANDROID_STAGE_URL=http://10.77.8.134:8080/stage
# #472: the Presenter Stage APK the watchdog installs on TVs that need our own
# app. Dev deploys to /opt/presenter-dev (the default path is /opt/presenter).
Environment=PRESENTER_ANDROID_STAGE_APK=/opt/presenter-dev/presenter-stage.apk
Environment=RUST_LOG=info,tower_http=debug
# Offset integration ports to avoid conflicts with production
Environment=PRESENTER_COMPANION_ENABLED=0
Environment=PRESENTER_COMPANION_PORT=28175
Environment=PRESENTER_OSC_LISTEN_PORT=49051

# Security hardening - no CAP_NET_BIND_SERVICE needed for port 8080
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/presenter-dev
PrivateTmp=true

# Resource limits — #335 (mirror of prod; see presenter.service for context).
CPUQuota=300%
MemoryMax=4G
TasksMax=512
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
